# scripts/ustc-grab.service
[Unit]
Description=USTC Grab Report Course Scheduler (main.py every 1min, update every 30min)
After=network.target
Wants=network.target

[Service]
Type=simple
[Service]
Type=simple
# 建议使用专用系统用户运行
User=ustc-grab
# 推荐部署到 /etc 目录
WorkingDirectory=/etc/ustc_grab_report_course
# 脚本路径
ExecStart=/etc/ustc_grab_report_course/scripts/ustc_grab_runner.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

# 安全与资源限制
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
# 读写路径 (需确保该用户对日志目录有写权限)
ReadWritePaths=/etc/ustc_grab_report_course /var/log/ustc_grab

[Install]
WantedBy=multi-user.target